# ChessGame/chess-csharp.service
[Unit]
Description=Chaos Chess C# Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Trading-bot/ChessGame
Environment=DOTNET_ROOT=/root/.dotnet
Environment=PATH=/root/.dotnet:/root/.dotnet/tools:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/Trading-bot/ChessGame/bin/Release/net8.0/ChessGame
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
